# installers/wallpaper-changer/wallpaper-changer.service
[Unit]
Description=Runs a dotfiles wallpaper changer script

[Service]
ExecStart=%h/.dotfiles/installers/wallpaper-changer/wallpaper-changer.sh
