# shell/profile.d/example.sh
# this file will be sourced in .profile uppon terminal start
# well, this exacly not, but all other ending with .sh 
