# shell/profile.d/read.me
# *.sh files in this dir will be sourced in .profile
