# installers/shell/profile.d/golang.sh
PATH=/home/jezek/.go/bin:$PATH
export GOPATH=/home/jezek/.go:/home/jezek/Programy/go
